# RPiTargetFiles/etc/supervisor/conf.d/PF_main_script.conf
[program:PF_main_script]
command=/usr/local/bin/PF_main.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/PF_main.err.log
stdout_logfile=/var/log/PF_main.out.log 
